# sigma-rules/proc_create_schtasks_persistence.yml
title: Scheduled Task Creation via Schtasks
id: a1b2c3d4-0003-4e5f-a003-argus00000003
status: test
description: Detects scheduled task creation using schtasks.exe with /create flag. Used for persistence in Argus stress test (T1053.005). Schtasks spawned as child of PowerShell or cmd.exe is particularly suspicious.
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://github.com/farrukhCTI/soc-homelab
author: Farrukh (farrukhCTI)
date: 2026-05-16
tags:
    - attack.persistence
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains: '/create'
    suspicious_parent:
        ParentImage|endswith:
            - '\powershell.exe'
            - '\cmd.exe'
            - '\wscript.exe'
            - '\cscript.exe'
            - '\mshta.exe'
    condition: selection and suspicious_parent
falsepositives:
    - Legitimate software scheduling tasks during installation
    - Windows Update and system maintenance tasks (use ParentImage to filter)
level: high
fields:
    - CommandLine
    - ParentImage
    - ParentCommandLine
    - User
